# etc/worker.service
[Unit]
Description=Worker Service
After=network.target

[Service]
ExecStart=/opt/worker/worker
Restart=always
RestartSec=10s

# User isolation
User=worker
Group=worker

# Minimal capability for cgroup management
AmbientCapabilities=CAP_SYS_ADMIN
NoNewPrivileges=yes

# Cgroup v2 delegation
Delegate=yes
DelegateControllers=cpu memory io pids

# Enable accounting
CPUAccounting=yes
MemoryAccounting=yes
IOAccounting=yes
TasksAccounting=yes

# Use dedicated slice
Slice=worker.slice

# Basic settings
WorkingDirectory=/opt/worker
LimitNOFILE=4096
StandardOutput=journal
StandardError=journal

# Setup controllers after service starts
ExecStartPost=/bin/bash -c 'sleep 1; echo "+cpu +memory +io +pids" > /sys/fs/cgroup/worker.slice/worker.service/cgroup.subtree_control || true'

# Cleanup job cgroups on stop
ExecStopPost=/bin/bash -c 'find /sys/fs/cgroup/worker.slice/worker.service -name "job-*" -type d -exec rmdir {} \; 2>/dev/null || true'

[Install]
WantedBy=multi-user.target