# debian/prerm
#!/bin/bash
set -e

case "$1" in
    remove|upgrade|deconfigure)
        # Stop the service
        if systemctl is-active --quiet joblet; then
            systemctl stop joblet
        fi

        # Disable the service
        if systemctl is-enabled --quiet joblet; then
            systemctl disable joblet
        fi
        ;;
esac

exit 0